# gtpd.socket
[Unit]
Description=Gtpd API Socket
PartOf=gtpd.service

[Socket]
ListenStream=/run/gtpd

[Install]
WantedBy=sockets.target
